# config/runit/wsdd2/run
#!/bin/sh -x
sleep 10
exec /usr/sbin/wsdd2
