# config/runit/wsdd2/run
#!/bin/sh
sleep 10
exec /usr/sbin/wsdd2
