# config/runit/avahi/run
#!/bin/sh -x
exec avahi-daemon --no-rlimits
