# config/runit/samba/run
#!/bin/sh -x
sleep 2
exec smbd --foreground
